# etc/systemd/system/nextcloud-sync.slice
[Unit]
Description=Nextcloud Sync Slice
DefaultDependencies=no
Before=slices.target

[Slice]
CPUQuota=30%
MemoryMax=2G
